# launchd/com.wizz-cmd.icloud-nfs-exporter.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wizz-cmd.icloud-nfs-exporter</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/HydrationDaemon</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/__USERNAME__/Library/Logs/icloud-nfs-exporter.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/__USERNAME__/Library/Logs/icloud-nfs-exporter.err</string>
    <key>WorkingDirectory</key>
    <string>/Users/__USERNAME__</string>
</dict>
</plist>
